# tests/client/rtclient-supervisor.conf
# Sample supervisor configuration for the client
# apt-get install nodejs supervisor
# copy (a modified version of) this config to /etc/supervisor/conf.d/
# supervisorctl reload
# for more help:
# supervisorctl help

[program:rtclient]
command=/usr/bin/nodejs /data/project/parsoid/js/tests/client/rtclient-cluster.js  /data/project/parsoid/js/tests/client/config.js

# Make sure that the dir *containing* the db file is writable to this user too
user=nobody

# Keep a log
redirect_stderr=true
stdout_logfile=/var/log/rtclient.log

# Automatically restart, no matter what the return status is
autorestart=true
